# crypto-scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cryptoturkey.scraper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "/Users/anadoluaydinli/Desktop/crypto website project/crypto-rewards-comparison-turkey" && source .venv/bin/activate && python backend/scraper_service.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/anadoluaydinli/Desktop/crypto website project/crypto-rewards-comparison-turkey/backend/logs/launcher-error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/anadoluaydinli/Desktop/crypto website project/crypto-rewards-comparison-turkey/backend/logs/launcher-output.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist> 